# shell/profile.d/ssh-ident
#!/bin/zsh

export BINARY_SSH=/usr/bin/ssh

chmod -R u=rwX,go= ~/.ssh

